# etc/supervisord.d/sensu-uchiwa.ini
[program:sensu-uchiwa]
command=uchiwa -c /etc/uchiwa/config.json
directory=/usr/local/src/github.com/sensu/uchiwa/
autostart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
